# config/app_supervisord.conf
[supervisord]
nodaemon=true
[program:php-fpm]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/sbin/php-fpm
user=web
autostart=true
autorestart=true
priority=5
